# etc/crontab
# this crontab will be reloaded by the script
# deploy/deploy_server.sh on deployment

15 4 * * * cd app && source .env/bin/activate && ENV=production ./manage.py clearsessions
